# initramfs/etc/inittab
tty2::respawn:/sbin/agetty -8 -s 38400 tty2 linux
